# app/config/services.d/apache/run
#!/usr/bin/execlineb -P
/bin/execenv
/usr/sbin/apachectl -f /etc/apache2/httpd.conf -DFOREGROUND